# playbooks/templates/scraper.timer.j2
[Unit]
Description=Run Job Scraper every x hours
Requires={{ service_name }}.service

[Timer]
Unit={{ service_name }}.service
OnBootSec=5min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
